# deployments/tgcc.service
[Unit]
Description=tgcc — Telegram-Claude Code Bridge
Documentation=https://github.com/jaekwon-park/tgcc
After=network.target

[Service]
Type=simple
User=insainty
WorkingDirectory=/home/insainty
EnvironmentFile=/opt/tgcc/bin/.env
# Place the tmux control socket outside PrivateTmp so spawned tmux/claude
# sessions survive a tgcc restart (KillMode=process) and the freshly started
# tgcc can re-attach to them. A dedicated dir (not /tmp) also keeps tgcc's tmux
# server isolated from the user's own tmux sessions.
Environment=TMUX_TMPDIR=/opt/tgcc/run
ExecStart=/opt/tgcc/bin/tgcc serve
Restart=always
RestartSec=5

# Graceful shutdown. KillMode=process so only the bot process is signalled on
# stop/restart — tmux and its child claude sessions keep running and are
# reconciled (re-attached) by the new tgcc instance instead of being SIGKILLed.
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=10

# Security hardening
# Note: ProtectHome is intentionally NOT set. Claude Code (spawned per session)
# needs to write to several paths under $HOME beyond ~/.claude — notably ~/.local/
# (Claude version cache), ~/.cache, and Node's npm/npx temp dirs. Enabling
# ProtectHome=read-only + an explicit ReadWritePaths list caused Claude's TUI to
# silently fail to render (process running, pane blank, no transcript ever
# written). ProtectSystem=strict + NoNewPrivileges still constrain the bot
# itself; tighten ProtectHome again only after auditing every path Claude Code
# touches at startup and per-spawn.
#
# PrivateTmp is intentionally OFF. With KillMode=process, tmux and its claude
# children survive a tgcc restart — but systemd tears down the old instance's
# private /tmp on restart, leaving the surviving claude processes pointing at a
# dead /tmp. Newly spawned claude then hangs during startup because it cannot
# extract its node native modules to /tmp. Sharing the host /tmp keeps it stable
# across restarts. /tmp must be in ReadWritePaths or ProtectSystem=strict would
# mount it read-only (breaking hooks and temp-file writes).
NoNewPrivileges=yes
PrivateTmp=no
ProtectSystem=strict
ReadWritePaths=/tmp /opt/tgcc/bin /opt/tgcc/run /home/insainty/projects /home/insainty/.claude /home/insainty/.claude.json

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tgcc

[Install]
WantedBy=default.target
